# debian/rtl8723ds-bt.service
[Unit]
Description=Realtek rtl8723ds bluetooth hciattach

[Service]
ExecStart=/sbin/start-rtl-bluetooth

[Install]
WantedBy=multi-user.target
